# crontab/laravel-schedule

* * * * * php /var/www/html/artisan schedule:run >> /var/www/html/storage/logs/crontab.log 2>&1
